# srcs/wordpress/supervisord.conf
[unix_http_server]
file=/run/supervisor/supervisor.sock

[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor/supervisor.sock
prompt=mysupervisor

[program:nginx]
command=nginx -g "daemon off;"
killasgroup=true
stopasgroup=true
redirect_stderr=true

[program:php-pfm7]
command=php-fpm7 -F
stdout_logfile=/var/log/php-fpm/stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/php-fpm/stderr.log
stderr_logfile_maxbytes=0
exitcodes=0
